# bin/start
### BEGIN INIT INFO
# Provides:          Start button listener
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Start button listener
# Description:       Enable service provided by daemon.
### END INIT INFO

#!/usr/bin/env python

import sys, os, time
import RPi.GPIO as GPIO

sys.path.append("/home/pi/UTR5/buttons/")

pid = str(os.getpid())
pidfile = "/tmp/robot-runit.pid"

if os.path.isfile(pidfile):
	print "%s already exists, exiting" % pidfile
	sys.exit()
file(pidfile, 'w').write(pid)

try:
	from start_button import listen_for_start
	print "RUNIT SERVICE LISTENING"
	listen_for_start()
	print "START BUTTON HAS BEEN PRESSED"
	time.sleep(60)
	print "END"
	GPIO.cleanup()

finally:
	os.unlink(pidfile)
