# bin/stop
### BEGIN INIT INFO
# Provides:          Stop button listener
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Stop button listener
# Description:       Enable service provided by daemon.
### END INIT INFO

#!/usr/bin/env python

import sys, os, time
import RPi.GPIO as GPIO

sys.path.append("/home/pi/UTR5/buttons/")

from stop_button import listen_for_stop

print "RUNIT STOP SERVICE LISTENING"

listen_for_stop()

print "STOP BUTTON HAS BEEN PRESSED"

time.sleep(60)

#import subprocess
#subprocess.call(os.path.expanduser("~/UTR5/bin/emergency-stop") #set and unset all GPIO pins and force kill the start thread with the kernel

print "END"

GPIO.cleanup()
